# _archive/ai-hub-effects.service
[Unit]
Description=AI Hub Effects v2 (Rembg + Upscale)
After=network.target

[Service]
Type=simple
User=pepe
WorkingDirectory=/mnt/seagate/ai-hub-gateway
ExecStart=/home/pepe/ai_env/bin/python effects_service_real.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target