# _archive/musetalk.service
[Unit]
Description=MuseTalk Lip-Sync Service
After=network.target

[Service]
Type=simple
User=pepe
WorkingDirectory=/mnt/seagate/MuseTalk
ExecStart=/home/pepe/ai_env/bin/python app.py --ip 0.0.0.0 --port 8041
Restart=always
RestartSec=10
Environment=PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True

[Install]
WantedBy=multi-user.target